# infra/systemd/telemux.service.template
[Unit]
Description=brainstack telemux Telegram/Codex control plane
After=network-online.target

[Service]
Type=simple
WorkingDirectory={{PRODUCT_REPO}}/apps/telemux
EnvironmentFile={{CONFIG_ROOT}}/telemux.runtime.env
EnvironmentFile={{CONFIG_ROOT}}/telemux.secrets.env
ExecStart={{BUN_BIN}} --no-env-file run {{PRODUCT_REPO}}/apps/telemux/src/main.ts
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
